# scripts/launchd/com.user.weekly-research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.en.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.weeklyresearch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/DEIN_BENUTZERNAME/Documents/Intelligence-Hub/tools/weekly_research.py</string>
    </array>

    <!-- Jeden Montag um 07:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>HIER_API_KEY_EINTRAGEN</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/DEIN_BENUTZERNAME/Documents/Intelligence-Hub</string>

    <key>StandardOutPath</key>
    <string>/Users/DEIN_BENUTZERNAME/Documents/Intelligence-Hub/outputs/notes/weekly_research_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/DEIN_BENUTZERNAME/Documents/Intelligence-Hub/outputs/notes/weekly_research_stderr.log</string>
</dict>
</plist>
