# scripts/launchd/com.user.clippersync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.en.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.clippersync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/DEIN_BENUTZERNAME/Documents/intelligence-hub/tools/clipper_sync.sh</string>
    </array>

    <!-- Reagiert sofort wenn eine neue Datei im iCloud-Clipper-Ordner landet -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/DEIN_BENUTZERNAME/Library/Mobile Documents/com~apple~CloudDocs/Clipper</string>
    </array>

    <!-- Beim Mac-Start automatisch laden -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logs für Debugging -->
    <key>StandardOutPath</key>
    <string>/Users/DEIN_BENUTZERNAME/Documents/intelligence-hub/outputs/notes/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/DEIN_BENUTZERNAME/Documents/intelligence-hub/outputs/notes/launchd_stderr.log</string>
</dict>
</plist>
